# deployment/supervisor-chat-websocket.conf
; Copy to: /etc/supervisor/conf.d/caf-chat-websocket.conf
; Then: sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start caf-chat-websocket
;
; This runs the standalone Node.js websocket server from:
;   /www/wwwroot/cafremarket.co.mz/chat-ws-node/server.js
;
; Verify paths:
;   which node
;   ls -la /www/wwwroot/cafremarket.co.mz/chat-ws-node

[program:caf-chat-websocket]
process_name=%(program_name)s
command=/usr/bin/node /www/wwwroot/cafremarket.co.mz/chat-ws-node/server.js
directory=/www/wwwroot/cafremarket.co.mz/chat-ws-node
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=false
stdout_logfile=/www/wwwroot/cafremarket.co.mz/storage/logs/chat-socket-node.log
stderr_logfile=/www/wwwroot/cafremarket.co.mz/storage/logs/chat-socket-node-error.log
stopwaitsecs=3600
